# ansible/roles/common/templates/user/profile.d/50-editor.sh.j2
{{ ansible_managed | comment }}
export EDITOR=/usr/bin/vim
