# ansible/roles/desktop/templates/user/profile.d/50-linode.sh.j2
{{ ansible_managed | comment }}
export LINODE_CLI_TOKEN={{ pzzl_linode_token }}

# if [ -x /usr/local/bin/linode-cli ] ; then
#     <( /usr/local/bin/linode-cli completion bash )
# fi
